# systemd/deezer-stay-loggedin.timer
[Unit]
Description=Extend Deezer login session

[Timer]
OnCalendar=daily
Persistant=true

[Install]
WantedBy=timers.target
